# telegram-scraper/ops/systemd/telegram-scraper.service
[Unit]
Description=Telegram IPTV refresh pipeline

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/aiostreams/telegram-scraper
ExecStartPre=/usr/bin/docker build -t telegram-scraper .
ExecStart=/home/ubuntu/aiostreams/scripts/run_iptv_refresh_pipeline.sh
